# install/zynthomania.service.in
[Unit]
Description=Zynthomania launcher
After=sound.target jackd.service zynaddsubfx.service

[Service]
User=USERNAME
Group=audio
ExecStart=COMMAND_LINE
CPUSchedulingPolicy=rr
CPUSchedulingPriority={user/sub_priority}

Restart=on-failure
RestartSec=2

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=Zynthomania

[Install]
WantedBy=multi-user.target
